# ansible/ziti-nonprod/roles/ziti/templates/healthcheck.service.j2
# =============================================================================
# Ziti Health Check Systemd Service
# =============================================================================
# Generated by Ansible - DO NOT EDIT MANUALLY
# Template: roles/ziti/templates/healthcheck.service.j2
#
# This provides a health check endpoint for ALB target health checks.
# =============================================================================

[Unit]
Description=Ziti Health Check Server
Documentation=https://openziti.github.io/
After=network-online.target

[Service]
Type=simple
User={{ ziti_user }}
Group={{ ziti_group }}

ExecStart=/usr/bin/python3 {{ ziti_home }}/bin/healthcheck.py {{ ziti_health_port }}

Restart=always
RestartSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths={{ ziti_home }}

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ziti-healthcheck

[Install]
WantedBy=multi-user.target

