# ansible/ziti-nonprod/roles/ziti/templates/healthcheck.service.j2
# =============================================================================
# Ziti Health Check Systemd Service
# =============================================================================
# Managed by Ansible - DO NOT EDIT MANUALLY
# Simple HTTP server for ALB health checks on port 8080
# =============================================================================

[Unit]
Description=Ziti Health Check Server
Documentation=https://openziti.github.io/
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory={{ ziti_healthcheck_dir }}

ExecStart=/usr/bin/python3 {{ ziti_healthcheck_dir }}/healthcheck.py

Restart=always
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=healthcheck

[Install]
WantedBy=multi-user.target
