# ansible/ziti-nonprod/roles/ziti/templates/ziti-router.service.j2
# =============================================================================
# Ziti Router Systemd Service
# =============================================================================
# Managed by Ansible - DO NOT EDIT MANUALLY
# =============================================================================

[Unit]
Description=Ziti Router
Documentation=https://openziti.github.io/
After=network-online.target ziti-controller.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory={{ ziti_router_dir }}

ExecStart={{ ziti_bin }} router run {{ ziti_router_dir }}/router.yaml

Restart=always
RestartSec=5
TimeoutStartSec=120

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ziti-router

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
