# ansible/ziti-nonprod/roles/ziti/templates/ziti-controller.service.j2
# =============================================================================
# Ziti Controller Systemd Service
# =============================================================================
# Generated by Ansible - DO NOT EDIT MANUALLY
# Template: roles/ziti/templates/ziti-controller.service.j2
# =============================================================================

[Unit]
Description=Ziti Controller
Documentation=https://openziti.github.io/
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ ziti_user }}
Group={{ ziti_group }}
WorkingDirectory={{ ziti_home }}/controller

ExecStart={{ ziti_home }}/bin/ziti controller run {{ ziti_home }}/controller/controller.yaml

Restart=always
RestartSec=5

# Allow binding to privileged ports (443)
AmbientCapabilities=CAP_NET_BIND_SERVICE

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths={{ ziti_home }}
ReadWritePaths=/var/log/ziti

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ziti-controller

[Install]
WantedBy=multi-user.target
