# ansible/ziti-nonprod/roles/ziti/templates/ziti-controller.service.j2
# =============================================================================
# Ziti Controller Systemd Service
# =============================================================================
# Managed by Ansible - DO NOT EDIT MANUALLY
# =============================================================================

[Unit]
Description=Ziti Controller
Documentation=https://openziti.github.io/
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory={{ ziti_controller_dir }}

ExecStart={{ ziti_bin }} controller run {{ ziti_controller_dir }}/controller.yaml

Restart=always
RestartSec=5
TimeoutStartSec=120

# Allow binding to privileged ports (443)
AmbientCapabilities=CAP_NET_BIND_SERVICE

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ziti-controller

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
